# debian/libeufin-nexus.target
[Unit]
Description=LibEuFin Nexus
After=postgres.service network.target

Wants=libeufin-nexus-ebics-fetch.service
Wants=libeufin-nexus-ebics-submit.service
Wants=libeufin-nexus-httpd.service

[Install]
WantedBy=multi-user.target